# chardev_v2/10-mysimple.rules
# Runs chmod a+rw for mycdev events.

KERNEL=="mycdev", RUN+="/bin/chmod a+rw /dev/mycdev"
